# scripts/monitor/com.intraday.forward_watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.intraday.forward_watchdog</string>

  <key>WorkingDirectory</key>
  <string>/Users/taeyeong/intraday_trading</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>scripts/monitor/forward_watchdog.py</string>
    <string>--archive-root</string>
    <string>archive</string>
    <string>--poll-interval-sec</string>
    <string>600</string>
    <string>--threshold-multiplier</string>
    <string>1.5</string>
    <string>--cooldown-sec</string>
    <string>21600</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/forward_watchdog.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/forward_watchdog.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SEAL_OPEN</key>
    <string>1</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
